# protected/extensions/yii-crontab/crontabs/reportmail_crontab
*/1 * * * * php D://workspace//fruit-growers-v2//cron//yiicmd.php reportemail